# jobs/zookeeper/templates/zkcleanup.cron
3 * * * *   root    /var/vcap/jobs/zookeeper/bin/zkcleanup.sh 1>> /var/vcap/data/sys/log/zookeeper/zkcleanup.stdout.log 2>> /var/vcap/data/sys/log/zookeeper/zkcleanup.stderr.log