# jobs/zookeeper/templates/zkcleanup.cron
3 * * * *   root    /var/vcap/jobs/zookeeper/bin/zkcleanup.sh > /var/vcap/data/sys/log/zookeeper/zkcleanup.stdout.log >> /var/vcap/data/sys/log/zookeeper/zkcleanup.stderr.log